# docker/stream-recorder/crontab
* * * * * appuser cvlc "http://media.kcrw.com/pls/kcrwsimulcast.pls" --sout "file/mp3:/audio/KCRW-$(date '+\%Y-\%m-\%d-\%H-\%M').mp3" --run-time=10 --stop-time=10 vlc://quit
